# thewall.service
[Unit]
Description=Keeps the white walkers in look

[Service]
ExecStart=/home/pedrohlc/.local/bin/bluefence.py

[Install]
WantedBy=default.target
